# docker/webtop/rootfs/etc/xdg/autostart/wx-download.desktop
[Desktop Entry]
Type=Application
Name=WX Download Services
Exec=/usr/local/bin/wx-session-start
Terminal=false
X-GNOME-Autostart-enabled=true

